# deploy/ansible/templates/brai-api.service.j2
[Unit]
Description=Brai API {{ item.key }}
Wants=brai-temporal.service
After=network.target brai-temporal.service

[Service]
Type=simple
WorkingDirectory={{ brai_env_root }}/{{ item.value.path }}/source/services/brai_api
{% if item.key == 'prod' %}
EnvironmentFile={{ brai_api_env_source }}
ExecStart=/usr/bin/env BRAI_ENVIRONMENT={{ item.value.goal_agent_queue_environment }} PORT={{ item.value.api_port }} BRAI_RELEASE_DIR={{ brai_repo }}/deploy/releases BRAI_INBOX_STORAGE_ROOT={{ brai_repo }}/data/inbox-attachments BRAI_CODEX_BIN=/srv/opt/codex-cli/bin/codex BRAI_CODEX_MODEL=gpt-5.4-mini BRAI_CODEX_TIMEOUT_MS=20000 CODEX_HOME=/srv/opt/codex-runtime/brai HOME=/srv/opt/codex-runtime/brai {{ brai_node_bin }} src/index.js
{% else %}
EnvironmentFile={{ brai_env_root }}/{{ item.value.path }}/brai-api.env
ExecStart=/usr/bin/env BRAI_ENVIRONMENT={{ item.value.goal_agent_queue_environment }} PORT={{ item.value.api_port }} BRAI_RELEASE_DIR={{ brai_repo }}/deploy/releases BRAI_INBOX_STORAGE_ROOT={{ brai_env_root }}/{{ item.value.path }}/data/inbox-attachments BRAI_CODEX_BIN=/srv/opt/codex-cli/bin/codex BRAI_CODEX_MODEL=gpt-5.4-mini BRAI_CODEX_TIMEOUT_MS=20000 CODEX_HOME=/srv/opt/codex-runtime/brai HOME=/srv/opt/codex-runtime/brai {{ brai_node_bin }} src/index.js
{% endif %}
Restart=always
RestartSec=3
TimeoutStopSec=45
User={{ brai_service_user }}
{% if item.key == 'prod' %}
Group={{ brai_service_group }}
{% else %}
Group={{ brai_deploy_user }}
{% endif %}
SupplementaryGroups={{ brai_service_group }} {{ brai_deploy_user }} {{ brai_source_group }} {{ brai_codex_exec_group }} {{ brai_codex_auth_group }}
UMask=0002
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
