# deploy/systemd/brai-socraticode-watcher.service
[Unit]
Description=Brai SocratiCode watcher
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mark
Group=mark
WorkingDirectory=/srv/projects/brai
Environment=BRAI_SOCRATICODE_ROOT=/srv/projects/brai
Environment=BRAI_SOCRATICODE_WATCHER_INTERVAL_MS=60000
ExecStart=/srv/opt/node-v22.16.0/bin/node /srv/projects/brai/scripts/brai-socraticode-watcher.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
