# admin/deploy/brai-admin.service
[Unit]
Description=Brai Admin
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/projects/brai/admin
EnvironmentFile=/etc/brai/brai-api.env
EnvironmentFile=/etc/brai/admin.env
ExecStart=/srv/opt/node-v22.16.0/bin/node node_modules/next/dist/bin/next start -H 127.0.0.1 -p 3040
Restart=always
RestartSec=3
User=mark
Group=mark
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
